# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
loglevel=info

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/app
environment=ENV="production"
stdout_logfile=/tmp/uvicorn.log
autorestart=true
stopsignal=INT

[program:whatsapp]
command=node index.js
directory=/app/whatsapp-service
environment=WA_PORT="3001",WA_WEBHOOK_URL="http://localhost:8000/api/v1/whatsapp/webhook"
stdout_logfile=/tmp/whatsapp.log
autorestart=true
stopsignal=INT
